# system/crontab
@hourly /Users/jheppler/acad/dissertation/wordcount.sh
