# system/crontab
0   *   *   *   *   /home/lmullen/acad/research/dissertation/wordcount.sh