# system/crontab
0   *   *   *   *   /home/jheppler/acad/research/dissertation/wordcount.sh